# linux/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:skynet_iot]
directory=/usr/iot
command=/usr/iot/skynet_iot/linux/x86_64/process-monitor -M 2 -D /usr/iot/skynet -S "sh /usr/iot/skynet_iot/startup.sh /usr/iot" /usr/iot/skynet/skynet iot/config
priority=1
autostart=true
autorestart=true
; user=frappe

[group:mqtt-clients]
programs=skynet_iot
