# linux/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:skynet_iot]
directory=/mnt
command=/mnt/skynet_iot/linux/x86_64/process-monitor -M 2 -D /mnt/skynet -S "sh /mnt/skynet_iot/startup.sh /mnt" /mnt/skynet/skynet iot/config
priority=1
autostart=true
autorestart=true
; user=frappe

[group:mqtt-clients]
programs=skynet_iot
